# launchd/com.gemini.telemetry.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.gemini.telemetry.agent</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/bretbarnard/Projects/gemini-mvp/scripts/run-telemetry-agent.sh</string>
      <string>30</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/bretbarnard/Projects/gemini-mvp</string>

    <key>StandardOutPath</key>
    <string>/Users/bretbarnard/Projects/gemini-mvp/var/telemetry_agent.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/bretbarnard/Projects/gemini-mvp/var/telemetry_agent.err</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>
  </dict>
</plist>
